# files/etc/uci-defaults/99-mwan3-track-host-routes
#!/bin/sh
# Enable host routes for mwan3track targets (nagual2 fork).
if ! uci -q get mwan3.globals.track_host_routes >/dev/null; then
	uci set mwan3.globals.track_host_routes='1'
	uci commit mwan3
fi
exit 0
